# launchd/com.knowledge.channel-telegram.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.knowledge.channel-telegram</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/knowledge/bin/kb-channel-telegram-run</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.bun/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>KB_HUB</key>
        <string>__HOME__/knowledge</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/knowledge/logs/channel-telegram.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/knowledge/logs/channel-telegram.stderr.log</string>
</dict>
</plist>
